# driver/99-bscpcie.rules
ACTION=="add",SUBSYSTEM=="pci",ATTR{vendor}=="0x10ee", ATTR{device}="0x7028", RUN+="/sbin/modprobe bdbmpcie"
KERNEL=="bdbm_regs0",MODE="666"
